# cool-apps.service
[Unit]
Description=Refresh cool-apps terminal cheat sheet
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cool-apps-refresh
StandardOutput=null
StandardError=journal
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
