# cool-apps.timer
[Unit]
Description=Weekly refresh of cool-apps terminal cheat sheet

[Timer]
OnCalendar=Sun *-*-* 10:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
